# acer_c720/xubuntu14.04/etc/xdg/autostart/compton.desktop
[Desktop Entry]
Type=Application
Terminal=false
Name=Compton
Comment=Compton compositor
Exec=compton -b
OnlyShowIn=XFCE;
